# scripts/appdaemon.service
[Unit]
Description=AppDaemon service for Home Assistant

[Service]
ExecStart=/usr/local/bin/appdaemon

[Install]
WantedBy=multi-user.target